# debian/postinst
systemctl daemon-reload
systemctl enable dnshield
systemctl start dnshield