# debian/prerm
systemctl stop dnshield
systemctl disable dnshield
systemctl daemon-reload

rm -f /etc/dnshield-conf.json